# crontab-open/install/etc/cron.d/crc

# made by tfq
@reboot username vcvlc #boot vlc
* 10 * * 1-5 username killall vlc && export DISPLAY=:1 && cvlc --loop --fullscreen --no-osd /crc/crontab-open/* #open
* 19 * * 1-5 username killall vlc && export DISPLAY=:1 && cvlc --loop --fullscreen --no-osd /crc/crontab-closed/* #close

